# systemd/jupyter.service
[Unit]
Description=Single user Jupyter server
After=network.target
#StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=5
User=simview
ExecStart=/opt/singularity/bin/singularity run --nv -B /media --app jupyter /home/simview/burkhard/src/singularity-containers/bh_dev.sif

[Install]
WantedBy=multi-user.target
